# backend/php/supervisor_api.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm -D

[program:composer]
command=/usr/local/bin/composer install
chdir=/var/www/html
autostart=true
autorestart=false

[program:laravel-jobs-queue]
command=php /var/www/html/artisan queue:work --queue=zoho --sleep=1 --tries=1
directory=/var/www/html
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/zoho-jobs-queues.log

[program:laravel]
command=php /var/www/html/artisan serve --host=0.0.0.0 --port=8000
directory=/var/www/html
autostart=true
autorestart=true
redirect_stderr=true

